# net/tollgate-rs/files/etc/init.d/tollgate-rs
#!/bin/sh /etc/rc.common
# procd init script for the TollGate-RS node daemon.

START=99
USE_PROCD=1

BIN=/usr/bin/tollgate
CONF=/etc/tollgate/tollgate.yaml

start_service() {
	local enabled listen

	config_load tollgate-rs
	config_get enabled  main enabled 1
	config_get listen   main listen

	[ "$enabled" = "0" ] && {
		echo "tollgate-rs: disabled in UCI (tollgate-rs.main.enabled=0)" >&2
		return 0
	}

	procd_open_instance
	procd_set_param command    "$BIN" --config "$CONF" serve
	[ -n "$listen" ] && procd_set_param command "$BIN" --config "$CONF" serve --listen "$listen"
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}

reload_service() {
	restart
}

service_triggers() {
	procd_add_reload_trigger tollgate-rs
}
